# MATLAB/SupportPackages/R2018b/toolbox/target/supportpackages/raspberrypi/src/mwrunonboot.service
[Unit]
Description=Simulink Run on boot
After=multi-user.target

[Service]
Type=idle
RemainAfterExit=yes
ExecStart=/usr/local/bin/MW_runSimulinkModel

[Install]
WantedBy=multi-user.target

